# bin/beam.service
# beam.service

[Unit]
Description=Beam Service
After=multi-user.target

[Service]
Type=simple
User=pi
Environment=PYTHONPATH=/home/pi/workspace/beam/
ExecStart=/home/pi/workspace/beam/venv/bin/python /home/pi/workspace/beam/beam/runner.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
